# com.note-sorter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.note-sorter</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Update these paths to match your system -->
        <string>/Users/YOUR_USERNAME/.local/bin/uv</string>
        <string>run</string>
        <string>--script</string>
        <string>/Users/YOUR_USERNAME/.local/bin/note-sorter</string>
    </array>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>WatchPaths</key>
    <array>
        <string>/Users/YOUR_USERNAME/obsidian-drop</string>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/.config/note-sorter/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/.config/note-sorter/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/YOUR_USERNAME/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
